# huntarr/rootfs/etc/services.d/huntarr/finish
#!/command/execlineb -S0
# ==============================================================================
# Home Assistant Add-on: Huntarr
# Take down the service
# ==============================================================================

s6-svscanctl -t /var/run/s6/services
