# huntarr/rootfs/etc/services.d/huntarr/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Start Huntarr service
# ==============================================================================

bashio::log.info "Starting Huntarr..."

cd /app || exit 1

# Run Huntarr
