# huntarr/rootfs/etc/services.d/huntarr/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Add-on: Huntarr
# Runs Huntarr
# ==============================================================================

bashio::log.info "Starting Huntarr..."

# Export environment variables
export TZ=$(bashio::config 'TZ' 'America/New_York')

# Start Huntarr
cd /app || exit 1
exec s6-setuidgid abc python3 main.py
