# huntarr/rootfs/etc/services.d/huntarr/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: Huntarr
# Runs Huntarr
# ==============================================================================

bashio::log.info "Starting Huntarr..."

# Export environment variables
export TZ=$(bashio::config 'TZ')

# Export custom environment variables
if bashio::config.has_value 'env_vars'; then
    bashio::log.info "Setting custom environment variables..."
    for var in $(bashio::config 'env_vars|keys'); do
        name=$(bashio::config "env_vars[${var}].name")
        value=$(bashio::config "env_vars[${var}].value")
        bashio::log.info "Setting ${name}"
        export "${name}=${value}"
    done
fi

# Start Huntarr
bashio::log.info "Starting Huntarr application..."
exec /app/huntarr
